# dot_config/containers/systemd/music/soulsync.container
[Unit]
Description=SoulSync Music Discovery
Wants=network-online.target
After=network-online.target slskd.service

[Container]
Image=docker.io/boulderbadgedad/soulsync:latest
AutoUpdate=registry
Network=music.network
PublishPort=8008:8008
# Spotify OAuth callback listener. Bound to localhost so it's only reachable
# on-host or via an SSH tunnel during the one-time Spotify authentication.
PublishPort=127.0.0.1:8888:8888

# Run as uid 0, which maps to the host user (turisa) under rootless podman,
# matching slskd/navidrome so all three write the shared library as the same
# host user. The entrypoint detects it's already at the target uid and skips
# the gosu privilege-drop.
Environment=PUID=0
Environment=PGID=0
Environment=UMASK=022
Environment=TZ=America/Sao_Paulo

Volume=soulsync-config:/app/config
Volume=soulsync-data:/app/data
Volume=soulsync-logs:/app/logs
Volume=/mnt/x9pro/downloads:/app/downloads
Volume=/mnt/x9pro/music:/app/Transfer

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
